# src/systemd/dtranslatebot.service.in
[Unit]
Description=Discord Translation Bot
After=network.target

[Service]
User=dtranslatebot
Group=dtranslatebot
WorkingDirectory=@dtranslatebot_SERVICE_WORKDIR@
ExecStart="@CMAKE_INSTALL_FULL_BINDIR@/dtranslatebot" "@CMAKE_INSTALL_FULL_SYSCONFDIR@/dtranslatebot.json"
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
